# scripts/drone-management/DRONE_artifacts.timer
[Unit]
Description=Manage DRONE artifacts timer

[Timer]
OnUnitActiveSec=4h
OnBootSec=1h

[Install]
WantedBy=timers.target
