# scripts/drone-management/DRONE_artifacts.service
[Unit]
Description=Manage DRONE artifacts

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/DRONE_manage_artifacts.sh
